# crontab
* * * * * /bin/bash -l -c 'cd /Users/beth/Projects/bloc/blocitoff && RAILS_ENV=production bundle exec rake todo:delete_items --silent'
